# deploy/supervisor.conf
; ── Daphne (serveur ASGI Django) ──────────────────────────────────────────────
[program:daphne]
command=/home/deploy/recrutement/venv/bin/daphne -b 127.0.0.1 -p 8001 recrutement.asgi:application
directory=/home/deploy/recrutement
user=deploy
autostart=true
autorestart=true
stdout_logfile=/home/deploy/recrutement/logs/daphne.log
stderr_logfile=/home/deploy/recrutement/logs/daphne_err.log
environment=DJANGO_SETTINGS_MODULE="recrutement.settings"

; ── Celery (worker taches asynchrones) ───────────────────────────────────────
[program:celery]
command=/home/deploy/recrutement/venv/bin/celery -A recrutement worker --loglevel=info
directory=/home/deploy/recrutement
user=deploy
autostart=true
autorestart=true
stdout_logfile=/home/deploy/recrutement/logs/celery.log
stderr_logfile=/home/deploy/recrutement/logs/celery_err.log
environment=DJANGO_SETTINGS_MODULE="recrutement.settings"

; ── Groupe ────────────────────────────────────────────────────────────────────
[group:recrutement]
programs=daphne,celery
